# src/orderdispatcher/crontab
* * * * * python3 /orderdispatcher.py >> /logfile