# src/orderstatusupdater/crontab
* * * * * python3 /orderstatusupdater.py >> /logfile